# docker/crond/crontab
# %M(0 - 59)
#   %H(0 - 23)
#     %d(1 - 31)
#       %m (1 - 12)
#         %w (0 - 6)
#           command
0 5 * * * /daily
0 6 * * 0 /weekly
